# ops/systemd/dsa110-autopipeline.timer
[Unit]
Description=Launch the DSA-110 batch pipeline daily at 02:00 UTC

[Timer]
OnCalendar=*-*-* 02:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
