# ops/systemd/dsa110-dashboard.service
[Unit]
Description=DSA-110 continuum QA dashboard (qa_server, port 8767)
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/data/dsa110-continuum-dashboard
Environment=PYTHONPATH=/data/dsa110-continuum-dashboard
Environment=DSA110_REPO_ROOT=/data/dsa110-continuum
EnvironmentFile=-/home/ubuntu/.config/dsa110/dashboard.env
ExecStart=/opt/miniforge/envs/casa6/bin/python -m uvicorn scripts.qa_server:app --host 0.0.0.0 --port 8767 --log-level warning
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
