# ops/systemd/dsa110-autopipeline.service
[Unit]
Description=DSA-110 continuum daily batch pipeline launcher

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/data/dsa110-continuum
Environment=PYTHONPATH=/data/dsa110-continuum
ExecStart=/opt/miniforge/envs/casa6/bin/python scripts/auto_pipeline.py
